# system-wake-up.service
[Unit]
Description=System Wake Up Action
[Service]
ExecStart=/usr/bin/true
